# playbook/stunnel-nfs-client.socket
[Unit]
Description=NFS over stunnel client

[Socket]
ListenStream={{ server_port }}
Accept=yes
TimeoutSec=60

[Install]
WantedBy=sockets.target
