# playbook/stunnel-nfs.socket
[Unit]
Description=NFS over STUNNEL

[Socket]
ListenStream={{ server_port }}
Accept=yes
TimeoutSec=60

[Install]
WantedBy=sockets.target
