# playbook/stunnel-nfs@.service
[Unit]
Description=NFS over STUNNEL Service

[Service]
ExecStart=/usr/bin/stunnel /etc/stunnel/stunnel-nfs.conf
StandardInput=socket
